# config/supervisord_content.conf
; ============================================================================
; supervisord_content.conf — Content VM process management
; AutoFarm Zero — Success Guru Network v6.0
;
; Install: cp config/supervisord_content.conf /etc/supervisor/conf.d/autofarm.conf
;          supervisorctl reread && supervisorctl update
;
; Content VM (A1 Flex: 3 OCPU, 20 GB RAM, 8 GB swap)
; Manages long-running daemons. Cron jobs are handled separately via crontab.
; ============================================================================

; ---------------------------------------------------------------------------
; Idle Guard Daemon
; Prevents OCI from reclaiming idle A1 instances by generating synthetic
; workload (SQLite ANALYZE, log compression) when CPU < 15% for >30 min.
; ---------------------------------------------------------------------------
[program:idle-guard]
command=/app/.venv/bin/python -m modules.infrastructure.idle_guard
directory=/app
user=autofarm
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/app/logs/idle_guard.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_redirect_stderr=true
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"

; ---------------------------------------------------------------------------
; Health Monitor Daemon
; Periodic system health checks — CPU, RAM, disk, DB connectivity, API health.
; Reports anomalies to Telegram and writes to health.log.
; ---------------------------------------------------------------------------
[program:health-monitor]
command=/app/.venv/bin/python -m modules.infrastructure.health_monitor
directory=/app
user=autofarm
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/app/logs/health.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_redirect_stderr=true
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"

; ---------------------------------------------------------------------------
; Group: all content-vm daemons
; ---------------------------------------------------------------------------
[group:autofarm-content]
programs=idle-guard,health-monitor
priority=999
