# config/supervisord_proxy.conf
; ============================================================================
; supervisord_proxy.conf — Proxy VM process management
; AutoFarm Zero — Success Guru Network v6.0
;
; Install: cp config/supervisord_proxy.conf /etc/supervisor/conf.d/autofarm-proxy.conf
;          supervisorctl reread && supervisorctl update
;
; Proxy VM (A1 Flex: 1 OCPU, 4 GB RAM)
; Manages the approval server, Telegram review bot, and token refresh daemon.
; Squid proxy instances are managed via systemd (squid@brand.service).
; ============================================================================

; ---------------------------------------------------------------------------
; Approval Server (Flask)
; HTTP server on port 8080 for review approve/reject actions.
; Endpoints: /approve/<token>, /reject/<token>, /health, /dashboard, /calendar
; ---------------------------------------------------------------------------
[program:approval-server]
command=/app/.venv/bin/python -m modules.review_gate.approval_server
directory=/app
user=autofarm
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/app/logs/approval_server.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_redirect_stderr=true
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",FLASK_ENV="production"

; ---------------------------------------------------------------------------
; Telegram Review Bot
; Long-running daemon that sends review packages to the Telegram review chat.
; Handles inline callback_data for approve/reject actions.
; Primary review channel — email is fallback only.
; ---------------------------------------------------------------------------
[program:telegram-review-bot]
command=/app/.venv/bin/python -m modules.notifications.telegram_bot
directory=/app
user=autofarm
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/app/logs/telegram_bot.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_redirect_stderr=true
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"

; ---------------------------------------------------------------------------
; Token Refresh Daemon
; Keeps OAuth tokens alive for all platform accounts.
; Refreshes tokens before publishing windows to prevent auth failures.
; ---------------------------------------------------------------------------
[program:token-refresher]
command=/app/.venv/bin/python -m account_manager.token_refresher
directory=/app
user=autofarm
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/app/logs/tokens.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_redirect_stderr=true
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"

; ---------------------------------------------------------------------------
; Group: all proxy-vm daemons
; ---------------------------------------------------------------------------
[group:autofarm-proxy]
programs=approval-server,telegram-review-bot,token-refresher
priority=999
